# adapter_startup.conf
# SNAPconnect - start a synapse adapter as a service
#
description "Start SNAPconnect"
start on runlevel [2]
stop on runlevel [!2]
respawn
exec python /home/snap/synapse_adapter.py
